# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command=sh -c "cd /etc/grafana-7.2.0/bin && ./grafana-server && kill -9 supervisord"
