# srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=sh -c "kill $(pidof sshd) || /usr/sbin/sshd -D && kill 1"
autostart=true
autorestart=true


[program:nginx]
command=sh -c "kill $(pidof nginx) || /usr/sbin/nginx -g 'daemon off;' && kill 1"
autostart=true
autorestart=true
