# srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=sh -c "/usr/sbin/sshd -D && kill -9 supervisord"

[program:nginx]
command=sh -c "/usr/sbin/nginx -g 'daemon off;' && kill -9 supervisord"
