# srcs/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=sh -c "kill $(pidof php-fpm7) || /usr/sbin/php-fpm7 -R -F && kill 1"
autostart=true
autorestart=true

[program:nginx]
command=sh -c "kill $(pidof nginx) || /usr/sbin/nginx -g 'daemon off;' && kill 1"
autostart=true
autorestart=true
